# deploy/ecoverde-backend.service
[Unit]
Description=EcoVerde Django API (Gunicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/ecoverde/backend
EnvironmentFile=/opt/ecoverde/.env
ExecStart=/opt/ecoverde/backend/.venv/bin/gunicorn ecoverde.wsgi:application --bind 127.0.0.1:8000 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
